# modules/cosmetic/templates/release.init.erb
#!/bin/bash
#
# release   This script updates the issue file.
#
# chkconfig: 2345 01 99
# description:  This updates the issue files using linux_logo to display \
#       lots of cool and useful extra info.

# Source function library.
. /etc/init.d/functions

# Release name
if [ -f /etc/fedora-release ]; then
    R=$(cat /etc/fedora-release)
elif [ -f /etc/redhat-release ]; then
    R=$(cat /etc/redhat-release | sed 's/ (.*)//g')
else
    R="Unknown distribution"
fi

# Kernel version - instead of #V since RHEL6 appended arch
V=$(basename $(uname -r) .$(uname -m))

# Format string to use
STRING="<%= @format_string %>"

# Check for linux_logo first
[ -x /usr/bin/linux_logo ] || exit 1

RETVAL=0

case "$1" in
    start|restart|reload|condrestart)
        echo -n $"Updating issue, issue.net and motd files: "
        /usr/bin/linux_logo -g -F "${STRING}" > /etc/issue
        RETVAL="$(($RETVAL+$?))"
        /usr/bin/linux_logo -g -F "${STRING}" -a > /etc/issue.net
        RETVAL="$(($RETVAL+$?))"
        /usr/bin/linux_logo -g -F "\n${STRING}" -a > /etc/motd
        RETVAL="$(($RETVAL+$?))"

        if [ "$RETVAL" -eq 0 ]; then
            action "" /bin/true
        else
            action "" /bin/false
        fi
        ;;
    status)
        cat /etc/motd
        RETVAL=0
        ;;
    stop)
        RETVAL=0
        ;;
    *)
        echo $"Usage: $0 {start|stop|restart|status}"
        RETVAL=1
        ;;
esac

exit $RETVAL

