# modules-wip/graylog2/templates/graylog2-web.init.erb
#!/bin/bash
#
# chkconfig: - 90 02
# description: Graylog2 Web Interface

### BEGIN INIT INFO
# Provides: graylog2-web
# Required-Start: $network $local_fs $remote_fs
# Required-Stop: $network $local_fs $remote_fs
# Should-Start: $syslog $named
# Should-Stop: $syslog $named
# Short-Description: start and stop graylog2-web
# Description: Graylog2 Web Interface.
### END INIT INFO

# Source function library.
. /etc/init.d/functions

# Source networking configuration.
. /etc/sysconfig/network

prog=graylog2-web
pidfile=/var/run/${prog}.pid
lockfile=/var/lock/subsys/${prog}

# Source local configuration.
. /etc/sysconfig/${prog}

start() {
    [ "$EUID" != "0" ] && exit 4
    [ "$NETWORKING" = "no" ] && exit 1
    [ -x /usr/bin/java ] || exit 5
    [ -f /etc/graylog2.conf ] || exit 6

    # Start daemons.
    echo -n $"Starting $prog: "
    if [ -f $lockfile ]; then
        echo -n "already running?"
        failure $"$prog already running?"
        echo
        exit 1
    fi
<% if @thin -%>
    /usr/bin/thin -u <%= @user %> -g <%= @user %> -e <%= @thin_environment %> -c /opt/graylog2-web -d -P ${pidfile} start
<% else -%>
    echo -n "unsupported web server type..."
    failure
    echo
    exit 1
<% end -%>
    RETVAL=$?
    [ $RETVAL -eq 0 ] && success || failure $"$prog start"
    echo
    [ $RETVAL -eq 0 ] && touch $lockfile
    return $RETVAL
}

stop() {
    [ "$EUID" != "0" ] && exit 4
    echo -n $"Shutting down $prog: "
    killproc -p ${pidfile} $prog
    RETVAL=$?
    echo
    rm -f $lockfile $pidfile
    return $RETVAL
}

# See how we were called.
case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  status)
    status -p ${pidfile} $prog
    ;;
  restart|force-reload)
    stop
    start
    ;;
  try-restart|condrestart)
    if status -p ${pidfile} $prog > /dev/null; then
        stop
        start
    fi
    ;;
  reload)
    exit 3
    ;;
  *)
    echo $"Usage: $0 {start|stop|status|restart|try-restart|force-reload}"
    exit 2
esac
